# root/etc/s6-overlay/s6-rc.d/jvm-args/up
s6-setuidgid abc /etc/s6-overlay/scripts/jvm-args
