# root/etc/s6-overlay/s6-rc.d/mc-custom-settings/up
/etc/s6-overlay/scripts/mc-custom-settings
