# root/etc/s6-overlay/s6-rc.d/mc-restart/up
# Main service name
define service_name "minecraft-server"

# If file exists
ifelse { s6-test -f "/tmp/${TMP_MC_RESTART}" }
{
    # Restart main service
    foreground { s6-echo "mc-restart: info: restarting service ${service_name}" }
    s6-svc -r "/run/service/${service_name}"

    # Remove tmp file
    s6-rmrf "/tmp/${TMP_MC_RESTART}"
}
