# root/etc/s6-overlay/s6-rc.d/mc-restart/up
# Main service name
define service_name "minecraft-server"

# Restart main service
foreground { s6-echo "s6-user: info: restarting service ${service_name}" }
s6-svc -r "/run/service/${service_name}"
