# root/etc/s6-overlay/s6-rc.d/accept-eula/up
#!/command/execlineb -P
with-contenv

# Get ACCEPT_EULA env variable value [true/false]
importas -D "" accept_eula ACCEPT_EULA

# If accept_eula is true
ifelse { s6-test "true" = ${accept_eula} }
{
    # Get DATA_DIR env variable
    importas -D "" data_dir DATA_DIR

    # If data_dir is not null
    if { s6-test -n ${data_dir} }

    # Set eula=true in data directory
    redirfd -w 1 "${data_dir}/eula.txt" s6-echo "eula=true"
}
